# system/nort.service
[Unit]
Description=Nort Analytics Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nort
WorkingDirectory=/opt/nort
ExecStart=/opt/nort/venv/bin/python run.py
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Environment
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/nort/data_store/device.json.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nort-agent

# Memory limits (L4-fix: leave 2 GB headroom for kernel/GPU on 8 GB Orin Nano)
MemoryMax=6G
MemoryHigh=5G
OOMScoreAdjust=-500

# Watchdog (L5-fix: 30 s; run.py sends sd_notify("WATCHDOG=1") every 10 s)
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
